# userspace/root/etc/sv/ncm-param-watcher/run
#!/bin/sh
# ncm-param-watcher - watches /data/params/d/UsbNcmEnabled and runs set_ncm.sh on value changes
exec 2>&1

WATCH_FILE="/data/params/d/UsbNcmEnabled"
WATCH_DIR=$(dirname "$WATCH_FILE")
STATE_FILE="/tmp/ncm-param-watcher.last"

# Get current value (empty string if file doesn't exist)
get_value() {
  cat "$WATCH_FILE" 2>/dev/null || echo ""
}

# Default NCM to enabled on first boot
if [ ! -f "$WATCH_FILE" ]; then
  mkdir -p "$WATCH_DIR"
  printf '1' > "$WATCH_FILE"
fi

# Run once at startup with current param value
/usr/comma/set_ncm.sh

get_value > "$STATE_FILE"

# Wait for the params directory to exist
while [ ! -d "$WATCH_DIR" ]; do
  sleep 2
done

# Watch the directory for changes to the param file
WATCH_FILENAME=$(basename "$WATCH_FILE")
exec inotifywait -m -e modify,create,delete "$WATCH_DIR" 2>/dev/null | while read -r dir event file; do
  [ "$file" = "$WATCH_FILENAME" ] || continue
  NEW_VALUE=$(get_value)
  LAST_VALUE=$(cat "$STATE_FILE" 2>/dev/null || echo "")
  if [ "$NEW_VALUE" != "$LAST_VALUE" ]; then
    /usr/comma/set_ncm.sh 2>/dev/null || true
    echo "$NEW_VALUE" > "$STATE_FILE"
  fi
done
