# userspace/root/etc/profile.d/gpu.sh
export RUSTICL_ENABLE=freedreno
export DISABLE_FAST_IDIV=1
export EMULATED_DTYPES=long
